# deploy/orcheeder.ini
[program:orcheeder]
command=/home/orcheeder/env/bin/python /home/orcheeder/fase2/manage.py run_gunicorn -b localhost:8081 -w 3 -t 300
directory=/home/orcheeder/fase2/
user=nobody
autostart=true
autorestart=true
redirect_stderr=True
